# scripts/launchd/dev.lantern.whatsapp-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>LANTERN_API_URL</key>
		<string>http://localhost:8080</string>
		<key>LANTERN_BRIDGE_PORT</key>
		<string>3100</string>
		<key>LOG_LEVEL</key>
		<string>info</string>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>
	<key>KeepAlive</key>
	<dict>
		<key>Crashed</key>
		<true/>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>Label</key>
	<string>dev.lantern.whatsapp-bridge</string>
	<key>ProgramArguments</key>
	<array>
		<string>__REPO_ROOT__/scripts/launchd/run-bridge-wrapper.sh</string>
		<string>__REPO_ROOT__/services/whatsapp-bridge</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>__HOME__/Library/Logs/Lantern/whatsapp-bridge.err.log</string>
	<key>StandardOutPath</key>
	<string>__HOME__/Library/Logs/Lantern/whatsapp-bridge.out.log</string>
	<key>ThrottleInterval</key>
	<integer>10</integer>
	<key>WorkingDirectory</key>
	<string>__REPO_ROOT__/services/whatsapp-bridge</string>
</dict>
</plist>
